# shaper
#!/bin/sh /etc/rc.common
# Mr Shaper - v2.1 (Subnet Aware)
# Reads subnet, mask, and rates from GPIO-Predefined [SHAPER] section
# Primary goal: cap video streaming to ~360-480p on metered 4G SIM
# Applied to eth1 (USB WAN / Mode1) only - proof of concept
# v2.1 Mar 2026 - Added CIDR support for /28 blocks

START=99
CONFIG="/root/GPIO-Predefined"

getcfg() {
    awk -F= -v sec="[$1]" -v key="$2" '
        $0==sec{f=1;next}
        /^\[/{f=0}
        f && /^#/{next}
        f && $1==key{print $2; exit}
    ' "$CONFIG"
}

start() {
    # Read shaper config
    SUB=$(getcfg SHAPER LIMITED_SUB)
    MASK=$(getcfg SHAPER SUB_MASK)
    RATE=$(getcfg SHAPER RATE)
    DEFAULT_RATE=$(getcfg SHAPER DEFAULT_RATE)

    # Fallback defaults
    SUB=${SUB:-192.168.231.80}
    MASK=${MASK:-28}
    RATE=${RATE:-6000kbit}
    DEFAULT_RATE=${DEFAULT_RATE:-8000kbit}

    IFACE="br-wifi_iso"  # Target the bridge directly

    # Clear old rules
    tc qdisc del dev $IFACE root 2>/dev/null

    # Root HTB qdisc
    tc qdisc add dev $IFACE root handle 1: htb default 10

    # Default class - full speed
    tc class add dev $IFACE parent 1: classid 1:10 htb rate $DEFAULT_RATE ceil $DEFAULT_RATE

    # Limited class - for target subnet
    tc class add dev $IFACE parent 1: classid 1:20 htb rate $RATE ceil $RATE

    # Attach fq_codel for stability
    tc qdisc add dev $IFACE parent 1:10 handle 10: fq_codel
    tc qdisc add dev $IFACE parent 1:20 handle 20: fq_codel

    # Apply filter using CIDR notation (Subnet/Mask)
    # Using 'ip src' because we are shaping traffic GOING TO the WAN from br-wifi_iso
    tc filter add dev $IFACE protocol ip parent 1: prio 1 u32 \
        match ip src ${SUB}/${MASK} flowid 1:20

    logger "Shaper: limiting ${SUB}/${MASK} to $RATE on $IFACE"
}

stop() {
    tc qdisc del dev br-wifi_iso root 2>/dev/null   #logic to stop br-wifi_iso interface    
                                                                                     
    logger "Shaper: stopped on br-wifi_iso interface"
}

